# settings/systemd/random-wallpaper.timer
[Unit]
Description=Randomize wallpaper every 5 minutes

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target