# settings/autostart/xfce4-clipman.desktop
[Desktop Entry]
Comment=Clipboard manager
Exec=xfce4-clipman
Icon=xfce4-clipman-plugin
Name=Clipman
Terminal=false
Type=Application
